# templates/service/master.service.erb
[Unit]
Description=<%= @masterServiceName %>
After=network.target


[Service]
<% @masterOptions.each do | option, value| -%>
# <%= "Environment=MESOS_#{option}=#{value}" %>
<%= "Environment=\"MESOS_#{option}=#{value}\"" %>
<% end -%>
Environment="MESOS_WORK_DIR=<%= @masterWorkDir %>"
Environment="MESOS_LOG_DIR=<%= @masterLogDir %>"
User=root
Group=root
UMask=0027
SyslogIdentifier=<%= @masterServiceName %>
WorkingDirectory=<%= @masterWorkDir %>
ExecStart=/usr/local/sbin/mesos-master
LimitNOFILE=16384
CPUAccounting=true
MemoryAccounting=true
Restart=always
RestartSec=20

[Install]
WantedBy=multi-user.target